# templates/etc/profile.d/miniconda.sh.j2
MINICONDA_DIR={{miniconda_parent_dir}}/miniconda
PATH=${MINICONDA_DIR}/bin:${PATH}
PYTHON=${MINICONDA_DIR}/bin/python


