# files/etc/profile.d/miniconda.sh
ANACONDA_DIR=/opt/miniconda
PATH=${ANACONDA_DIR}/bin:${PATH}
PYTHON=${ANACONDA_DIR}/bin/python


